# rotate-factorio.timer
[Unit]
Description=Weekly factorio save rotation

[Timer]
OnCalendar=Tue *-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target